# luci-app-pihole/root/etc/uci-defaults/10-luci-pihole
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@pihole[-1]
	add ucitrack pihole
	set ucitrack.@pihole[-1].init=pihole-FTL
	commit ucitrack
EOF

exit 0
